# scripts/network-namespace.service
# network-namespace.service - Network Namespace 开机启动服务
# 文件位置: /etc/systemd/system/network-namespace.service
#
# 安装方法:
#   sudo cp network-namespace.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable network-namespace.service
#
# 使用方法:
#   sudo systemctl start network-namespace.service   # 启动
#   sudo systemctl stop network-namespace.service    # 停止（会恢复原配置）
#   sudo systemctl status network-namespace.service  # 查看状态

[Unit]
Description=Network Namespace Setup for Agent isolation
After=network-online.target
Wants=network-online.target
Before=agent-eth1-ns.service agent-eth2-ns.service agent-eth3-ns.service agent-eth5-ns.service agent-eth6-ns.service agent-eth7-ns.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/SFW-TEST3/scripts/network-namespace-setup.sh setup-all
ExecStop=/opt/SFW-TEST3/scripts/network-namespace-setup.sh restore-all
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target